# s6_overlay/etc/services.d/s6log/run
#!/usr/bin/with-contenv bash
set -e

#
#   This service passes the S6-overlay services' logfile to stdout
#

exec tail -f -n 0 /opt/iobroker/log/s6-overlay/current
